# config/systemd/system/pi-web-restart.service
[Unit]
Description=Restart pi-web

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart pi-web
